# root/etc/services.d/proxy/run
#!/usr/bin/with-contenv sh

exec stolon-proxy --listen-address 0.0.0.0 --cluster-name stolon-cluster --store-backend consul --store-endpoints=consul:8500;
